# agent_supervisord.conf
[program:way-home]
command=/usr/bin/node /home/user/way-home/dist/index.js client
directory=/home/user/way-home
user=user