# server_supervisord.conf
[program:way-home]
command=way-home server
directory=/home/user/way-home
user=user